# services/receiver.service
[Unit]
Description=GNSS Data Receiver Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/kirill/Practice/receiver_service.py <receiver_name>
Restart=always
User=kirill
WorkingDirectory=/home/kirill/Practice

[Install]
WantedBy=multi-user.target
